# etc/cron.d/aprs_beacon
# APRS beacon
# m h day(month) month day(week) user command
*/5 * * * *     mmdvm  [ -x /home/mmdvm/aprs_beacon ] && /home/mmdvm/aprs_beacon
